# packaging/runtara-core/postinst
#!/bin/bash
set -e

# Create runtara user and group if they don't exist
if ! getent group runtara > /dev/null 2>&1; then
    groupadd --system runtara
fi

if ! getent passwd runtara > /dev/null 2>&1; then
    useradd --system --gid runtara --no-create-home --shell /usr/sbin/nologin runtara
fi

# Create config directory if it doesn't exist
mkdir -p /etc/runtara
chown root:runtara /etc/runtara
chmod 750 /etc/runtara

# Set permissions on config file
if [ -f /etc/runtara/core.conf ]; then
    chown root:runtara /etc/runtara/core.conf
    chmod 640 /etc/runtara/core.conf
fi

# Reload systemd and enable service
systemctl daemon-reload
systemctl enable runtara-core.service

echo "Runtara Core installed successfully."
echo "Edit /etc/runtara/core.conf to configure the service."
echo "Start with: systemctl start runtara-core"
